# systemd/eliza-voice-llm.service
[Unit]
Description=Eliza Voice LLM Endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/private/eliza-inference
ExecStart=%h/private/eliza-inference/scripts/start voice-llm --profile voice-gemma4-e4b-default
ExecStop=%h/private/eliza-inference/scripts/stop voice-llm --profile voice-gemma4-e4b-default
Restart=on-failure

[Install]
WantedBy=default.target
