# systemd/eliza-stt.service
[Unit]
Description=Eliza STT Endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/private/eliza-inference
ExecStart=%h/private/eliza-inference/scripts/start stt --profile stt-faster-whisper-large-v3
ExecStop=%h/private/eliza-inference/scripts/stop stt --profile stt-faster-whisper-large-v3
Restart=on-failure

[Install]
WantedBy=default.target
