# systemd/eliza-medium.service
[Unit]
Description=Eliza Medium Reasoning Model Endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/private/eliza-inference
ExecStart=%h/private/eliza-inference/scripts/start eliza-medium --profile eliza-medium-qwen-llamacpp-32k
ExecStop=%h/private/eliza-inference/scripts/stop eliza-medium --profile eliza-medium-qwen-llamacpp-32k
Restart=on-failure

[Install]
WantedBy=default.target
