# systemd/eliza-qwen-coder.service
[Unit]
Description=Eliza Qwen Coder Inference Endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
WorkingDirectory=%h/private/eliza-inference
ExecStart=%h/private/eliza-inference/scripts/start qwen-coder --profile qwen-shared-200k
ExecStop=%h/private/eliza-inference/scripts/stop qwen-coder --profile qwen-shared-200k
Restart=on-failure

[Install]
WantedBy=default.target
